# web/kiran-sau-problem/conf-files/cron.conf
[program:cron]
user=root
command=cron -f
nodaemon=false
autostart=true
autorestart=true
stderr_logfile=/var/log/python.err.log
stdout_logfile=/var/log/python.out.log
